# warwar 0.3.5 jan23 + lvv 1.11.4/lv/lv-link/etc/systemd/system/lv-link@.service
[Unit]
Description=lv-link
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/bin/lv-link-%i

[Install]
WantedBy=multi-user.target
