# warwar 0.3.2 aug9 + vision 1.11.1/lv/lv-link/etc/systemd/system/lv-keep-alive.service
[Unit]
Description=lv-link air keep alive service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/bin/lv-keep-alive.sh

[Install]
WantedBy=multi-user.target
